# scripts/sentinel-firstboot.service
[Unit]
Description=Sentinel First Boot Provisioning
After=network-online.target
Wants=network-online.target
Before=sentinel.service

[Service]
Type=oneshot
ExecStart=/home/sentinellite/sentinel/scripts/first-boot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
